# src/modules/225_awstats/awstats_monthly.cron
#!/bin/bash

for cfg in `find /etc/awstats -name 'awstats.*.conf' -printf '%f\n' | sed 's/^awstats\.\(.*\)\.conf/\1/'`
do
	/opt/ics/gotypo/awstats_update.sh $cfg fullyear >/dev/null
done
